# matchbox/assets/flatcar/oem-cloudinit.service
[Service]
Type=oneshot
ExecStart=/usr/bin/flatcar-cloudinit --oem=ec2-compat

[Install]
WantedBy=multi-user.target
